# root/etc/services.d/mpd/run
#!/usr/bin/with-contenv sh

# wait for services to come up
sleep 2
echo "Starting MPD Server..."
exec mpd --no-daemon --stdout -v